# files/etc/rc.local
#!/bin/bash
/etc/init.d/startupUbuntuRaoul
exit 0